# aliwal-modules/initscript/templates/ruby-init.erb
#! /bin/bash

### BEGIN INIT INFO
# Provides:          <%= @app %>-<%= @daemon %>-<%= @environment %>
# Required-Start:    $local_fs $remote_fs $network $syslog
# Required-Stop:     $local_fs $remote_fs $network $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: starts the daemon <%= @daemon %>
# Description:       starts the daemon <%= @daemon %> for app <%= @app %> at <%= @environment %> environment
### END INIT INFO

NAME="daemon $DAEMON_NAME"
DESC="Daemon $DAEMON_NAME for app $APP"
USER="<%= @user %>"
USER_PATH="<%= @home %>"
APP="<%= @app %>"
APP_PATH="$USER_PATH/$APP"
CD_TO_APP_DIR="cd $APP_PATH"
RBENV_PATH="$USER_PATH/.rbenv"
DAEMON_NAME=<%= @daemon %>
DAEMON="$APP_PATH/bin/$DAEMON_NAME"
DAEMON_OPTS=""
START_DAEMON_PROCESS="$DAEMON $DAEMON_OPTS"
PID=/var/run/$APP.pid

case "$1" in
  start)
        echo -n "Starting $DESC: "
        if [ `whoami` = root ]; then
          su - $USER -c "$CD_TO_APP_DIR > /dev/null 2>&1 && $START_DAEMON_PROCESS &"
        else
          $CD_TO_APP_DIR > /dev/null 2>&1 && $START_DAEMON_PROCESS &
        fi
        echo "$NAME."
        ;;
  stop)
        echo -n "Stopping $DESC: "
        kill -9 `cat $PID`
        echo "$NAME."
        ;;
  restart)
        echo -n "Restarting $DESC: "
        if [ -f "$PID" ]; then
          $0 stop
          sleep 1
        fi
        $0 start
        ;;
  reload)
        echo -n "Reloading $DESC configuration: "
        kill -HUP `cat $PID`
        echo "$NAME."
        ;;

  rotate)
        echo -n "Rotating logs $DESC: "
        kill -USR1 `cat $PID`
        echo "$NAME."
        ;;

  *)
        echo "Usage: $NAME {start|stop|restart|reload}" >&2
        exit 1
        ;;
esac

exit 0
